# systemd/user/mintupdate-automation-notifier.timer
[Unit]
Description=Update Manager Automatic Updates Notifier

[Timer]
OnStartupSec=10s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
